# auto-theme/aether-scheduler.service.d/auto-theme.conf
[Service]
ExecStartPost=/home/ioio/.local/bin/omarchy-auto-theme
